# systemd/assistant-web.service
[Unit]
Description=Personal Research Assistant Web UI
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/musa/Desktop/asistant
EnvironmentFile=/home/musa/Desktop/asistant/.env
ExecStart=/home/musa/Desktop/asistant/.venv/bin/python scripts/run_web.py
Restart=on-failure
RestartSec=5
User=musa
Group=musa

[Install]
WantedBy=multi-user.target
