# systemd/assistant-bot.service
[Unit]
Description=Personal Research Assistant Telegram Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/Desktop/asistant
EnvironmentFile=%h/Desktop/asistant/.env
ExecStart=%h/Desktop/asistant/.venv/bin/python scripts/run_bot.py
Restart=on-failure
RestartSec=5
User=musa
Group=musa

[Install]
WantedBy=multi-user.target
